# shani-deploy/deploy-image.service
[Unit]
Description=Shani OS Deployment Service
Documentation=man:deploy-image(1)
After=network-online.target
Requires=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/deploy-image
StandardOutput=journal
StandardError=journal
TimeoutStopSec=300
Restart=on-failure
RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target
