# shani-deploy/shani-rollback.service
[Unit]
Description=Rollback System on Boot Failure for Shani OS
Documentation=man:rollback-system(1)
Conflicts=rescue.service emergency.service
After=systemd-boot-failed.service
Before=shutdown.target
ConditionPathExists=/var/run/boot_failure

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rollback-system.sh
RemainAfterExit=true
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300
Restart=no

[Install]
WantedBy=multi-user.target

